# deploy/lm-database-api.service
[Unit]
Description=lm-database API server
After=network.target mongod.service
Wants=mongod.service

[Service]
Type=simple
User=lm-database
Group=lm-database
WorkingDirectory=/opt/lm-database
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/local/cargo/bin
EnvironmentFile=/opt/lm-database/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true
ProtectHome=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=
SystemCallFilter=@system-service
UMask=0077

# Allow writes only to required directories
ReadWritePaths=/var/lib/lm-database/temp /opt/lm-database/logs

[Install]
WantedBy=multi-user.target
